# systemd/netmon-ping.service
[Unit]
Description=NetMon Ping Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/netmon/scripts/ping_monitor.py
EnvironmentFile=/opt/netmon/config/secrets.env
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netmon-ping

[Install]
WantedBy=multi-user.target
